# systemd/.config/systemd/user/mping.timer
[Unit]
Description=Daily mping timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=mping.service

[Install]
WantedBy=timers.target
