# systemd/.config/systemd/user/wbg.service.d/override.conf
[Service]
ExecStart=/usr/sbin/wbg -s /home/morethan/Pictures/wallpapers/longway.png
